# tools/pkgs/01-php-setup/top/etc/cron.d/epi-php
MAILTO=infra@epiconcept.fr

9,39 * * * *	root	/usr/bin/php_clean_sessions

# EOF
