# php/pkgs/02-mod/top/etc/cron.d/epi-php
MAILTO=infra@epiconcept.fr

9,39 * * * *    root	/usr/local/bin/php_clean_sessions.sh

# EOF